# templates/instance/wireguard/wg-watch.timer
[Unit]
Description=Run WG Broker guardian every 5 mins

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
